# infra/systemd/gunicorn@.socket
[Unit]
Description=%i gunicorn socket

[Socket]
ListenStream=/srv/hr/run/%i.sock
SocketUser=hr
SocketMode=600

[Install]
WantedBy=sockets.target

